# data/storaged.service.in
[Unit]
Description=LVM Disk Manager
Documentation=man:udisks-lvm(8)

[Service]
Type=dbus
BusName=com.redhat.lvm2
ExecStart=@udisksdprivdir@/udisksd-lvm
