# data/storaged.service.in
[Unit]
Description=LVM Disk Manager
Documentation=man:storaged(8)

[Service]
Type=dbus
BusName=com.redhat.storaged
ExecStart=@storagedprivdir@/storaged
